# crontab
0 * * * * khll /usr/bin/python3 /path/to/PGC.py
